# 16-asus-xonar-u1.rules
ACTION=="add", SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0b05", ATTRS{idProduct}=="1743", RUN+="/usr/local/bin/xonard $root/$name"
